# templates/caddy.conf
description "Caddy Server startup script"

start on runlevel [2345]
stop on runlevel [016]

setuid caddy
setgid caddy

respawn
respawn limit 10 5

script
    exec /usr/bin/caddy
end script
